# templates/config/systemd/user/iwallpaper_linux.service
[Unit]
Description=iWallpaper: An intelligent wallpaper manager

[Service]
ExecStart=/home/kai/mypy/bin/iwallpaper

[Install]
WantedBy=default.target
